# docker/cron/cron.d/app
0 0 */15 * * cd /app && python app/cron.py >> /var/log/cron/app.log 2>&1
